# .config/systemd/user/borgmatic.timer.d/twice-a-week.conf
# This override for borgmatic timer changes execution frequency
# from daily backups at midnight to twice-a-week backups at 2AM.

[Timer]
OnCalendar=
OnCalendar=Mon,Thu 02:00:00
