# .config/systemd/user/mpd.socket.d/custom-port.conf
# This snippet ensures that the MPD server listens on preferred custom port.

[Socket]
ListenStream=
ListenStream=50012
